# legacy_archive/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=logs/supervisord.log
pidfile=/tmp/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:ingest_worker]
command=python ingest_history.py --mode=incremental
autostart=true
autorestart=true
startretries=3
stdout_logfile=logs/ingest_worker.log
stderr_logfile=logs/ingest_worker.err

[program:ingest_backfill]
command=python ingest_history.py --mode=deep_backfill
autostart=false
autorestart=false
stdout_logfile=logs/ingest_backfill.log
stderr_logfile=logs/ingest_backfill.err

[program:feature_refresher]
command=python feature_refresher.py --interval=30
autostart=true
autorestart=true
stdout_logfile=logs/feature_refresher.log

[program:live_engine]
command=python main.py
autostart=false         ; only start after verify OK
autorestart=true
startretries=5
stdout_logfile=logs/live_engine.log
stderr_logfile=logs/live_engine.err
